# Operation/exampleSystemd.service
[Unit]
Description=My Binary Service
After=network.target

[Service]
ExecStart=/path/to/your/binary
WorkingDirectory=/path/to/your/binary_directory
Restart=always
User=myuser
Group=mygroup
StandardOutput=append:/var/log/mybinary.log
StandardError=append:/var/log/mybinary.log

[Install]
WantedBy=multi-user.target
